# autoReverseSSH.conf
# autossh
# place into /etc/init/
description "autossh daemon"

#start on net-device-up IFACE=eth0
start on net-device-up IFACE=wlo1
stop on runlevel [01S6]

respawn
respawn limit 5 60 # respawn max 5 times in 60 seconds

script
	export AUTOSSH_PIDFILE=/var/run/autossh.pid
	export AUTOSSH_POLL=60
	export AUTOSSH_FIRST_POLL=30
	export AUTOSSH_GATETIME=0
	export AUTOSSH_DEBUG=1
	autossh -M 10003 -N -R 10002:localhost:22 remoteUser@remoteIPAddress -o "ServerAliveInterval 60" -o "ServerAliveCountMax 3" -i /home/USER/.ssh/privateRSAKey
end script


# and then log into the remoteUser@remoteIPAddress and run
# ssh remoteUser@localhost -p 10002

# (although autossh tends to just pick some arbitrary port) , 
# so look at the ps aux | grep autossh on the local machine to find which port it binds to
